# tquic-manager/scripts/tquic-manager.service
[Unit]
Description=TQUIC Manager Web Interface
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tquic-manager/backend
Environment=TQUIC_MANAGER_BIND=127.0.0.1
Environment=TQUIC_MANAGER_PORT=5000
ExecStart=/usr/bin/python3 /opt/tquic-manager/backend/app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
